# docker/app/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
socketfile=/var/run/supervisor.sock

[program:php-fpm]
command=php-fpm -F --nodaemonize
autostart=true
autorestart=true
user=application
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.out.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
user=application
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log 

[program:consumer-jira-event]
command=php /opt/app/bin/console.php messenger:consume jira_event_transport
autostart=true
autorestart=true
user=application
stderr_logfile=/var/log/supervisor/jira-event-consumer.err.log
stdout_logfile=/var/log/supervisor/jira-event-consumer.out.log 
